# sml-snmp-agent/files/sml_server.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2015 OpenWrt.org

START=77

start() {
	sml_server -i /dev/ttyUSB0
}

stop() {
	killall sml_server
}
